# archive/udoo-ref/UDOObuntu/lib/systemd/system/udoo-firstrun.service
[Unit]
Description=UDOO firstrun script

[Service]
ExecStart=/usr/sbin/udoo-firstrun
SyslogIdentifier=udoo-firstrun

[Install]
WantedBy=local-fs.target

