# archive/udoo-ref/UDOObuntu/lib/systemd/system/udooneo-bluetooth.service
[Unit]
Description=UDOO Neo Bluetooth Stack

[Service]
Type=exec
ExecStart=/usr/sbin/uimlauncher

[Install]
WantedBy=multi-user.target

