# archive/udoo-ref/UDOObuntu/var/lib/dpkg/info/locate.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/locate
/usr
/usr/bin
/usr/bin/locate.findutils
/usr/bin/updatedb.findutils
/usr/lib
/usr/lib/arm-linux-gnueabihf
/usr/lib/arm-linux-gnueabihf/locate
/usr/lib/arm-linux-gnueabihf/locate/frcode
/usr/share
/usr/share/doc
/usr/share/doc/locate
/usr/share/doc/locate/README.Debian
/usr/share/doc/locate/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/locate.findutils.1.gz
/usr/share/man/man1/updatedb.findutils.1.gz
/usr/share/man/man5
/usr/share/man/man5/locatedb.5.gz
/var
/var/cache
/var/cache/locate
/usr/share/doc/locate/changelog.Debian.gz
